# deploy/paddington.service
# Paddington always-on voice assistant (systemd USER service).
#
# Install on the Pi (user "lab", repo at ~/Labie, venv at ~/Labie/.venv):
#   mkdir -p ~/.config/systemd/user
#   cp deploy/paddington.service ~/.config/systemd/user/paddington.service
#   sudo loginctl enable-linger "$USER"        # run without an active login / across reboots
#   systemctl --user daemon-reload
#   systemctl --user enable --now paddington
#
# Manage:
#   systemctl --user status paddington
#   systemctl --user restart paddington
#   systemctl --user stop paddington           # frees the mic for the cry monitor
#   tail -f ~/paddington.log
#
# NOTE: the mic can only be owned by one program at a time, so this service and
# the cry-detection `listen`/`analyze` loop must not run on the same mic at once.

[Unit]
Description=Paddington voice assistant (Lullaby)
After=default.target sound.target

[Service]
Type=simple
WorkingDirectory=/home/lab/Labie
Environment=SDL_AUDIODRIVER=pulse
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/lab/Labie/.venv/bin/python -m lullaby --config /home/lab/Labie/config/demo.toml listen-assistant
StandardOutput=append:/home/lab/paddington.log
StandardError=append:/home/lab/paddington.log
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
